# .pkg-build.otfVRA/packaging/rhel/patch-agent-heartbeat.service
[Unit]
Description=PatchMaster Agent Heartbeat
After=network.target
Wants=patch-agent.service

[Service]
Type=simple
User=patchagent
Group=patchagent
EnvironmentFile=-/etc/patch-agent/env
WorkingDirectory=/opt/patch-agent
ExecStart=/opt/patch-agent/run-heartbeat.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
