# packaging/rhel/patch-agent.service
[Unit]
Description=PatchMaster Agent API
After=network.target

[Service]
Type=simple
User=patchagent
Group=patchagent
EnvironmentFile=-/etc/patch-agent/env
WorkingDirectory=/opt/patch-agent
ExecStart=/opt/patch-agent/run-api.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
